# winchester-parse/com.cloudflare.cw-parse.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.cloudflare.cw-parse</string>
    <key>ProgramArguments</key>
    <array>
        <string>/opt/homebrew/bin/cloudflared</string>
        <string>tunnel</string>
        <string>--config</string>
        <string>/Users/thorsnode/.cloudflared/cw-parse-config.yml</string>
        <string>run</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/Users/thorsnode/Library/Logs/cloudflared-cw-parse.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/thorsnode/Library/Logs/cloudflared-cw-parse.err.log</string>
</dict>
</plist>
